# scripts/hermes-dashboard-proxy.service
[Unit]
Description=Hermes Dashboard Host-Header Proxy
After=network-online.target hermes-dashboard.service
Wants=network-online.target

[Service]
User=hermes
Group=hermes
WorkingDirectory=/home/hermes
Environment=HOME=/home/hermes
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 /opt/hermes/scripts/run-hermes-dashboard-proxy.py
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
