# scripts/hermes-cron-tick@.service
[Unit]
Description=Hermes cron ticker for profile %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hermes
Group=hermes
WorkingDirectory=/home/hermes/work/%i
ExecStart=/opt/hermes/scripts/run-profile-cron-tick.sh %i
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
